# forwarding_bot.service
[Unit]
Description=Telegram Forwarding Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_LINUX_USER
WorkingDirectory=/path/to/forwarding_bot_gemini
EnvironmentFile=/path/to/forwarding_bot_gemini/.env
ExecStart=/path/to/forwarding_bot_gemini/venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forwarding_bot

[Install]
WantedBy=multi-user.target
